# edge/systemd/auxin-watchdog.service
[Unit]
Description=Auxin Automata Safety Watchdog (independent torque monitor)
After=network-online.target
Wants=network-online.target

# Watchdog runs independently — it MUST NOT depend on auxin-bridge.service.
# If the bridge crashes, the watchdog still halts the arm.

[Service]
Type=simple
User=jetson
Environment="ROS_DOMAIN_ID=0"
ExecStart=/bin/bash -c 'source /opt/ros/humble/setup.bash && source /home/jetson/auxin_ws/install/setup.bash && ros2 run auxin_edge safety_watchdog_node'
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
